# mnemosyne/stacks/jobiris/jobiris-monitor-daily.service
[Unit]
Description=JobIris job search monitor (daily profiles, via Docker)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=youruser
WorkingDirectory=/home/youruser/homelab-infra/mnemosyne/stacks/jobiris
ExecStart=/usr/bin/docker compose run --rm jobiris-monitor python3 job-monitor.py --schedule daily
