# mnemosyne/systemd/container-update-metrics.service
[Unit]
Description=Docker Container Update Metrics (Prometheus Textfile Collector)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/container-update-metrics.sh
# Write stdout/stderr to journald — view with: journalctl -u container-update-metrics
StandardOutput=journal
StandardError=journal
