# mnemosyne/systemd/fan-metrics.timer
[Unit]
Description=Run fan-metrics.sh every 30 seconds

[Timer]
OnBootSec=15s
OnUnitActiveSec=30s
AccuracySec=1s

[Install]
WantedBy=timers.target
